# systemd/ztpbootstrap-nginx.container
[Unit]
Description=ZTP Bootstrap Nginx Container

[Container]
ContainerName=ztpbootstrap-nginx
Image=docker.io/nginx:alpine
Pod=ztpbootstrap.pod

# Expose ports directly
PublishPort=443:443
PublishPort=80:80

# Volume mounts
Volume=/opt/containerdata/ztpbootstrap:/usr/share/nginx/html:ro
Volume=/opt/containerdata/ztpbootstrap/nginx.conf:/etc/nginx/conf.d/default.conf:ro
Volume=/opt/containerdata/certs/wild:/etc/nginx/ssl:ro

Environment="TZ=UTC"

# Health check
HealthCmd=["sh", "-c", "pgrep nginx > /dev/null || exit 1"]
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=on-failure

[Install]
WantedBy=ztpbootstrap-pod.service
